# vm-contestant/usr/lib/systemd/system/contestant-reset.service
[Unit]
Description=Contestant reset
[Service]
Type=oneshot
ExecStart=/usr/local/sbin/contestant-reset
[Install]
WantedBy=multi-user.target

