# vm-contestant/usr/lib/systemd/system/oiproctor_mac.service
[Unit]
Description=OI Proctor random MAC
[Service]
Type=oneshot
ExecStart=/usr/local/sbin/oiproctor_mac
[Install]
WantedBy=multi-user.target

