# vm-contestant/etc/systemd/system/oiproctor_mac.service
[Unit]
Description=Random OI Proctor MAC
[Service]
Type=oneshot
ExecStart=/usr/local/sbin/oiproctor_mac
[Install]
WantedBy=multi-user.target

